# roles/means_media/templates/data-mount.service.j2
[Unit]
Description=Data Mount

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/udisksctl mount -b /dev/mapper/{{ item.vg }}-{{ item.lv }}
ExecStop=/usr/bin/udisksctl unmount -b /dev/mapper/{{ item.vg }}-{{ item.lv }}
SuccessExitStatus=1

[Install]
WantedBy=default.target
